# packages/core-base/core_base-1.0.4.tar.gz/core_base-1.0.4/compose/supervisor/master-celery.conf
[program:celery_master]
autorestart=True
autostart=True
redirect_stderr=True
startsecs=10
satrtretries=3
command=celery -A sobase beat --loglevel=info --scheduler django_celery_beat.schedulers:DatabaseScheduler -f logs/celery.log
directory=/hbyapp/sobase/
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10
stdout_logfile = /var/log/supervisor/celery_worker_stdout.log
stopasgroup=true
killasgroup=true