# packages/core-base/core_base-1.0.4.tar.gz/core_base-1.0.4/compose/supervisor/websocket-base.conf
[program:celery_worker]
autorestart=True
autostart=True
redirect_stderr=True
startsecs=10
satrtretries=3
command=daphne -b 0.0.0.0 -p 9852 sobase.asgi:application
directory=/hbyapp/sobase/
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 10
stdout_logfile = /var/log/supervisor/celery_worker_stdout.log
stopasgroup=true
killasgroup=true